# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory: /var/www/
command=php artisan queue:work
